# examples/gunicorn/gunicorn.conf.example
description "Chishop Gunicorn server"
author "Rafael Durán Castañeda <rafael@bvox.net>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -c "/home/vagrant/chishop/bin/django-gunicorn.sh" vagrant
